# app-crypt/scdtools/files/scdrand.service
[Unit]
Description=Scdrand obtains random data from a smartcard and adds them to the kernel entropy pool.
After=local-fs.target

[Service]
Type=simple
Environment=INTERVAL=10
Environment=THRESHOLD=0
Environment=ENTROPY=8
Environment=BYTES=128
ExecStart=/usr/bin/scdrand -l -i ${INTERVAL} -t ${THRESHOLD} -e ${ENTROPY} ${BYTES}

[Install]
WantedBy=default.target
